# provisioning/playbooks/roles/ipfs/templates/ipfs.service.j2
[Unit]
Description=IPFS Daemon


[Service]
User="{{ipfs_user}}"
Group="{{ipfs_group}}"
ExecStart=/usr/local/bin/ipfs daemon


[Install]
WantedBy=multi-user.target
